# tools/templates/etc/init.d/clonemac
#!/bin/sh /etc/rc.common
# Copyright (C) 2014 OpenWrt.org

START=94
STOP=15

start() {
	ifconfig eth1 down
	ifconfig eth1 hw ether __clonemac__
	ifconfig eth1 up
}

stop() {
	echo "Stop."
}
